# infra/systemd/podman-platform-api.service
[Unit]
Description=PodMan platform API
After=network-online.target mongod.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/podman
Environment=NODE_ENV=production
EnvironmentFile=/root/podman/backend/.env
ExecStart=/usr/bin/node backend/dist/server.js
Restart=always
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
