# infra/systemd/podman-gemini-starter-agent.service
[Unit]
Description=PodMan LiveKit Gemini starter agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/podman/examples/livekit-gemini-hacker-starter/agent
Environment=PATH=/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/root/.local/bin/uv run agent.py dev
Restart=always
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
